# branches/remodex/src/tools/rc.d/init-script
##!/bin/bash
### BEGIN INIT INFO
# Provides:          remod: lightweight sauerbraten server
# Required-Start: $local_fs $network
# Required-Stop: $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: This provides a simple to use control for remod servers.
### END INIT INFO

NAME="Remod"
DESC="Remod Extended Sauerbraten Server"
DIR=/opt/remod
DAEMON=./remod
USERID=remod
pidfile=$DIR/$NAME.pid
running=false

process_running() {
	if [ -f $pidfile ] && [ "$( ps -p `cat $pidfile` -o comm= )" = 'remod' ]
	then
		return 0
	else
		return 1
	fi
}

startserver() {
	cd $DIR
	if process_running
	then
		echo "Server is already running! Please stop it first!"
		exit 0
	fi

	echo -n "Starting $DESC: "
	sudo -u $USERID bash <<EAT

	nohup $DAEMON 1>/dev/null 2>/dev/null &
        touch $pidfile
	jobs -p %- > $pidfile

	# Wait for the server to start completely
	until running=true 
	do
		if "$( ps -p `cat $pidfile` -o comm= )" = 'remod'
		then
			running=true
		else
			sleep 1
		fi
	done
EAT

	if process_running
	then
		echo -e '\e[1;32m[ OK ]\e[0m'
	else
		echo -e '\e[0;31m[fail]\e[0m'
	fi
}

stopserver() {
	if process_running
	then
		echo -n "Stopping $DESC: "
		sudo -u $USERID kill `cat $pidfile`
		until ! process_running
		do
			echo -n '.'
			sleep 1
		done
		echo -e "\e[1;32m[ OK ]\e[0m"
	else
		echo "Server is currently not running."
		exit 0
	fi
}

edit() {
	sudo -u $USERID nano "$DIR/server-init.cfg"
	if process_running
	then
		read -p "Shall I restart the server '$NAME'? [Y/n]: "
		case "$REPLY" in
	  	  n)
			exit 0
			;;
		  N)
			exit 0
			;;
		  *)
                stopserver
                startserver
			;;
		esac
	else
		exit
	fi
}

case "$1" in
  start)
	startserver
	;;
  stop)
	stopserver
	;;
  restart)
	stopserver
	startserver
	;;
  edit)
	edit
	;;
  status)
	if process_running
	then
		echo "Server currently running at PID `cat $pidfile`"
	else
		echo "Server is not running at the moment."
	fi
	;;
  *)
	echo -e "\e[1;33mUsage: \e[1;35m$0 \e[1;32m{start|stop|restart|status|edit}\e[00m" >&2
	exit 1
	;;
esac

# END OF SCRIPT


